# debian/zorp.service.in
[Unit]
Description=@ZORP_PRODUCT_NAME@

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/zorpctl start
ExecStop=/usr/sbin/zorpctl stop
ExecReload=/usr/sbin/zorpctl reload

[Install]
WantedBy=multi-user.target
