# debian/zorp.service.in
[Unit]
Description=The Zorp firewall starter

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/zorpctl start

[Install]
WantedBy=multi-user.target
